# net-irc/irker/files/irker.init
#!/sbin/runscript
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="/usr/bin/irkerd"
command_args="${IRKERD_OPTS}"
command_background="true"
pidfile="/var/run/${SVCNAME}.pid"

depend() {
	use net
}
